# webhook.service
[Unit]
Description=Telegram Bot Webhook Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bot_telegram
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/uvicorn webhook_app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

